# systemd/sgd.service
[Unit]
Description=agent-safe-guard native daemon
Requires=sgd.socket
After=sgd.socket

[Service]
Type=notify
ExecStart=/usr/local/bin/sgd
Restart=on-failure
RestartSec=1s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.claude %t/agent-safe-guard

[Install]
WantedBy=multi-user.target
